# systemd/hook-serve.service
[Unit]
Description=Hook Serve (Rust)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=relay
Group=relay
EnvironmentFile=/etc/hook-serve/relay.env
ExecStart=/usr/local/bin/hook-serve
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/
ReadWritePaths=/run

[Install]
WantedBy=multi-user.target
