# firecracker/systemd/external-chisel-check.service
[Unit]
Description=External Chisel Deep Check (run on separate host)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
Environment=FIRECRACKER_REPO_ROOT=/opt/hook-serve
EnvironmentFile=-/etc/firecracker/runtime.env
EnvironmentFile=-/etc/firecracker/alerts.env
EnvironmentFile=-/etc/firecracker/chisel-check.env
ExecStart=/usr/bin/env bash -lc 'exec "${FIRECRACKER_REPO_ROOT}/firecracker/watchdog/external-chisel-check.sh"'
