# firecracker/systemd/firecracker-overwatcher.service
[Unit]
Description=Firecracker Overwatcher (Kill Unresponsive MicroVM APIs)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
Environment=FIRECRACKER_REPO_ROOT=/opt/webhook-relay
EnvironmentFile=-/etc/firecracker/runtime.env
EnvironmentFile=-/etc/firecracker/overwatcher.env
ExecStart=/usr/bin/env bash -lc 'exec "${FIRECRACKER_REPO_ROOT}/firecracker/runtime/overwatcher.sh"'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
