# firecracker/systemd/firecracker@.service
[Unit]
Description=Firecracker MicroVM (%i)
After=network-online.target firecracker-network.service
Wants=network-online.target
Requires=firecracker-network.service

[Service]
Type=simple
User=root
Group=root
Environment=FIRECRACKER_REPO_ROOT=/opt/webhook-relay
EnvironmentFile=-/etc/firecracker/runtime.env
EnvironmentFile=/etc/firecracker/%i.env
ExecStart=/usr/bin/env bash -lc 'exec "${FIRECRACKER_REPO_ROOT}/firecracker/runtime/launch.sh" "%i" "${FIRECRACKER_CONFIG}" "${FIRECRACKER_SOCKET:-}"'
ExecStopPost=/usr/bin/env bash -lc 'exec "${FIRECRACKER_REPO_ROOT}/firecracker/runtime/cleanup.sh" "%i" "${FIRECRACKER_SOCKET:-}"'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
