# firecracker/systemd/firecracker-watchdog.service
[Unit]
Description=Firecracker Watchdog (auto-recovery + heartbeat)
After=network-online.target firecracker-network.service
Wants=network-online.target
Requires=firecracker-network.service

[Service]
Type=oneshot
User=root
Group=root
Environment=FIRECRACKER_REPO_ROOT=/opt/hook-serve
EnvironmentFile=-/etc/firecracker/runtime.env
EnvironmentFile=-/etc/firecracker/watchdog.env
EnvironmentFile=-/etc/firecracker/alerts.env
ExecStart=/usr/bin/env bash -lc 'exec "${FIRECRACKER_REPO_ROOT}/firecracker/watchdog/watchdog.sh"'
