# firecracker/systemd/pstore-collect.service
[Unit]
Description=Collect pstore kernel crash logs at boot
DefaultDependencies=no
After=local-fs.target
Before=multi-user.target
ConditionPathExists=/sys/fs/pstore

[Service]
Type=oneshot
User=root
Group=root
Environment=FIRECRACKER_REPO_ROOT=/opt/hook-serve
EnvironmentFile=-/etc/firecracker/runtime.env
EnvironmentFile=-/etc/firecracker/pstore-collect.env
ExecStart=/usr/bin/env bash -lc 'exec "${FIRECRACKER_REPO_ROOT}/firecracker/watchdog/pstore-collect.sh"'

[Install]
WantedBy=multi-user.target
